# ops/ai.openclaw.star-office-presence.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.openclaw.star-office-presence</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/obsidianassistant/clawd/apps/Star-Office-UI/ops/update-office-presence.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/obsidianassistant/clawd/apps/Star-Office-UI</string>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>60</integer>

    <key>StandardOutPath</key>
    <string>/Users/obsidianassistant/clawd/apps/Star-Office-UI/ops/update-office-presence.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/obsidianassistant/clawd/apps/Star-Office-UI/ops/update-office-presence.stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
  </dict>
</plist>
